# server/bt-reconnect.service
[Unit]
Description=Bluetooth reconnect watchdog for PT-P300BT
After=bluetooth.target

[Service]
Type=oneshot
ExecStart=/home/kalle/bt-reconnect.sh
