# server/label-printserver.service
[Unit]
Description=Label Print Server (Flask API + mobile webapp, gunicorn)
After=network.target bluetooth.target rfcomm-ptp300bt.service
Wants=rfcomm-ptp300bt.service

[Service]
Type=simple
User=kalle
WorkingDirectory=/home/kalle
ExecStart=/home/kalle/printserver-env/bin/gunicorn --workers 2 --bind 0.0.0.0:8080 --access-logfile - --error-logfile - print_server:app
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/kalle/logs

[Install]
WantedBy=multi-user.target
